# etc/cron.d/42-lgcron
@reboot lg sleep 30; tmux new-session -d -s screensaver "/home/lg/bin/screensaver.py /home/lg/bin/tour.sh"
