# etc/cron.d/42-lgcron
@reboot lg tmux new -ds screensaver "/home/lg/bin/screensaver.py /home/lg/bin/tour.sh"
